# target/linux/qualcommax/ipq807x/base-files/etc/init.d/arista-c360-radio-mode
#!/bin/sh /etc/rc.common

START=01

boot()
{
	/usr/sbin/arista-c360-radio-mode apply-early
}

start()
{
	/usr/sbin/arista-c360-radio-mode apply-early
}
