# target/linux/qualcommax/ipq807x/base-files/etc/uci-defaults/99_arista-c360-wireless-layout
#!/bin/sh

[ "$(board_name)" = "arista,c360" ] || exit 0

# Normalize AP-C360 radio order and mode-specific device settings only.
/usr/sbin/arista-c360-radio-mode apply-wireless

exit 0
